# contrib/crawl-livegames.service
Description=Crawl Live Games API

[Service]
ExecStart=/app/live-games/contrib/run.sh

[Install]
WantedBy=multi-user.target
